# pi/rc.local
#!/bin/sh

# keep net on because systemd/networkmanager is apparently garbage
while : ; do
    wpa_supplicant -i eth0 -c /etc/wpa_supplicant.conf -D wired &

    sleep 5

    dhclient eth0

    sleep 10

    # dns.watch
    echo 'nameserver 84.200.69.80' >/etc/resolv.conf

    sleep 30

    # uni net likes to keep allowing ping but randomly firewall me...
    while curl eff.org >/dev/null ; do
        sleep 30
    done
done &

exit 0
